# proxy.service
[Unit]
Description=Go Reverse Proxy for PR Previews
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
# このパスはプロキシを配置する場所に合わせて変更してください
WorkingDirectory=/opt/wol-proxy
ExecStart=/opt/wol-proxy/proxy-server -config /opt/wol-proxy/routes.json -port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
